# service/httpd/run
#!/bin/sh
exec httpd -c ./httpd.conf -p 0.0.0.0:8080
