# service/althttpd/run
#!/bin/sh
exec althttpd --user nobody --popup --root /www --port 8080
